# data/foldyd.service.in
[Unit]
Description=Folder Manager Dameon
After=default.target

[Service]
ExecStart=@bindir@/foldyd
Restart=always

[Install]
WantedBy=default.target
